# deploy/systemd/wa-monitor.service
[Unit]
Description=WA Monitor (whatsapp-web.js)
After=network-online.target
Wants=network-online.target

[Service]
Type=simple
WorkingDirectory=/opt/wa_2
ExecStart=/usr/bin/npm start
Restart=always
RestartSec=3
Environment=PORT=3000
Environment=HOST=127.0.0.1

#
# By default the service runs as root (systemd default).
# If you want to run as a non-root user later, add:
# User=wa
# Group=wa

[Install]
WantedBy=multi-user.target
